# pizero/hdmioff.timer
[Unit]
Description=Run hdmioff 10 minutes after boot

[Timer]
OnBootSec=10min

[Install]
WantedBy=timers.target
